# scripts/systemd-crm.service
[Unit]
Description=CRM Imobiliaria (FastAPI)
After=network.target

[Service]
User=www-data
WorkingDirectory=/var/www/crm/backend
EnvironmentFile=/var/www/crm/.env
ExecStart=/var/www/crm/backend/venv/bin/gunicorn -k uvicorn.workers.UvicornWorker app.main:app -b 127.0.0.1:8000 --workers 2 --timeout 120
Restart=always

[Install]
WantedBy=multi-user.target
